# deploy/systemd/agent-mockingbird.service
[Unit]
Description=Agent Mockingbird API and Dashboard
After=network.target executor.service
Wants=network.target executor.service

[Service]
Type=simple
User=agent-mockingbird
Group=agent-mockingbird
WorkingDirectory=/srv/agent-mockingbird/app
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=AGENT_MOCKINGBIRD_HOST=0.0.0.0
Environment=AGENT_MOCKINGBIRD_CONFIG_PATH=/var/lib/agent-mockingbird/agent-mockingbird.config.json
Environment=AGENT_MOCKINGBIRD_DB_PATH=/var/lib/agent-mockingbird/agent-mockingbird.db
Environment=AGENT_MOCKINGBIRD_EXECUTOR_ENABLED=true
Environment=AGENT_MOCKINGBIRD_EXECUTOR_BASE_URL=http://127.0.0.1:8788
Environment=AGENT_MOCKINGBIRD_EXECUTOR_WORKSPACE_DIR=/var/lib/agent-mockingbird/executor-workspace
Environment=AGENT_MOCKINGBIRD_EXECUTOR_DATA_DIR=/var/lib/agent-mockingbird/executor
Environment=AGENT_MOCKINGBIRD_EXECUTOR_UI_MOUNT_PATH=/executor
Environment=AGENT_MOCKINGBIRD_EXECUTOR_HEALTHCHECK_PATH=/executor
Environment=AGENT_MOCKINGBIRD_EXECUTOR_MODE=embedded-patched
Environment=AGENT_MOCKINGBIRD_MEMORY_API_BASE_URL=http://127.0.0.1:3001
Environment=AGENT_MOCKINGBIRD_MEMORY_WORKSPACE_DIR=/srv/agent-mockingbird/app
Environment=OPENCODE_CONFIG_DIR=__AGENT_MOCKINGBIRD_OPENCODE_CONFIG_DIR__
Environment=OPENCODE_DISABLE_PROJECT_CONFIG=1
Environment=OPENCODE_DISABLE_EXTERNAL_SKILLS=1
ExecStart=/usr/bin/env bun apps/server/src/index.ts
Restart=always
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/srv/agent-mockingbird/app /var/lib/agent-mockingbird

[Install]
WantedBy=multi-user.target
