# deploy/systemd/opencode.service
[Unit]
Description=OpenCode Sidecar for Agent Mockingbird
After=network.target
Wants=network.target

[Service]
Type=simple
User=agent-mockingbird
Group=agent-mockingbird
WorkingDirectory=/srv/agent-mockingbird/app
Environment=AGENT_MOCKINGBIRD_PORT=3001
Environment=AGENT_MOCKINGBIRD_MEMORY_API_BASE_URL=http://127.0.0.1:3001
ExecStart=/usr/bin/env opencode serve --hostname 127.0.0.1 --port 4096 --print-logs --log-level INFO
Restart=always
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/srv/agent-mockingbird/app /var/lib/agent-mockingbird

[Install]
WantedBy=multi-user.target
