# deploy/systemd/wafflebot.service
[Unit]
Description=Wafflebot API and Dashboard
After=network.target opencode.service
Wants=network.target opencode.service

[Service]
Type=simple
User=wafflebot
Group=wafflebot
WorkingDirectory=/srv/wafflebot/app
Environment=NODE_ENV=production
Environment=PORT=3001
Environment=WAFFLEBOT_CONFIG_PATH=/var/lib/wafflebot/wafflebot.config.json
Environment=WAFFLEBOT_DB_PATH=/var/lib/wafflebot/wafflebot.db
Environment=WAFFLEBOT_MEMORY_WORKSPACE_DIR=/srv/wafflebot/app
ExecStart=/usr/bin/env bun src/index.ts
Restart=always
RestartSec=2
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/srv/wafflebot/app /var/lib/wafflebot

[Install]
WantedBy=multi-user.target
